# stockapp-scheduler.service
[Unit]
Description=Stock Analyzer Scheduler Service
# Ensure network is up before starting
After=network.target

[Service]
# IMPORTANT: Replace 'hasher' if your user is different
User=hasher
Group=hasher
# Run from within the backend directory
WorkingDirectory=/home/hasher/Stock_Analysis/backend
# Load environment variables from .env file in the *project root*
EnvironmentFile=/home/hasher/Stock_Analysis/.env
# Set PYTHONPATH so imports work when run from backend dir
# Execute scheduler.py using the full path to python3 (adjust if needed)
ExecStart=/bin/bash -c 'cd /home/hasher/Stock_Analysis/backend && PYTHONPATH=/home/hasher/Stock_Analysis/backend /usr/bin/python3 scheduler.py'
Restart=always
RestartSec=10s # Restart after 10 seconds if it crashes

[Install]
WantedBy=multi-user.target
