# stockapp-web.service
[Unit]
Description=Gunicorn instance to serve Stock Analyzer Flask app
After=network.target # Ensure network is up

[Service]
# IMPORTANT: Replace 'hasher' if your user is different
User=hasher
Group=hasher
# Run from within the backend directory
WorkingDirectory=/home/hasher/Stock_Analysis/backend
# Load environment variables from .env file in the *project root*
EnvironmentFile=/home/hasher/Stock_Analysis/.env
# Execute gunicorn pointing to app:app within the backend directory
# Using absolute path for socket file
ExecStart=/home/hasher/Stock_Analysis/backend/venv/bin/gunicorn --workers 3 --bind unix:/home/hasher/Stock_Analysis/stockapp.sock -m 007 app:app
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
